# pkg/hqplayerd.service
[Unit]
Description=HQPlayer Embedded daemon
Requires=network-online.target sound.target
After=network-online.target sound.target local-fs.target remote-fs.target

[Service]
Type=simple
ExecStart=/usr/bin/hqplayerd
Restart=always
RestartSec=2
EnvironmentFile=-/etc/default/hqplayerd
Nice=-10
IOSchedulingClass=realtime
LimitMEMLOCK=4G
LimitNICE=-10
LimitRTPRIO=98

[Install]
WantedBy=multi-user.target
