# pkg/usb-affinity.service
[Unit]
Description=Bind USB Interrupts to CPU 5
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'for i in $(grep xhci_hcd /proc/interrupts | awk -F: "{print \$1}" | tr -d " "); do echo 20 > /proc/irq/$i/smp_affinity; done'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
