# pkg/usb-affinity.service
[Unit]
Description=Bind USB Interrupts to CPU 5
After=systemd-modules-load.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'for i in $(awk -F: "/xhci_hcd/ {gsub(/ /,\"\",\$1); print \$1}" /proc/interrupts); do echo 20 > /proc/irq/$i/smp_affinity; done'

[Install]
WantedBy=multi-user.target
